# ansible/k3s/roles/security/files/k3s.service.d/admission-dependency.conf
[Unit]
Requires=admission-controller.service
After=admission-controller.service
