# ansible/k3s/roles/security/files/k3s.service.d/dependencies.conf
[Unit]
Requires=k3s-config-init.service admission-controller.service
After=k3s-config-init.service admission-controller.service
